# etc/supervisor/conf.d/mysqld.conf
[program:mysqld]
command=/scripts/mysql/mysql-start.sh
